# lotus-exporter-farcaster/lotus-exporter-farcaster.cron
# EXPORT LOTUS TO NODE-EXPORTER / PROMETHEUS
* * * * *   {{ LOTUS_USER }} if /usr/local/bin/lotus-exporter-farcaster.py > {{ NODE_EXPORTER_PATH }}/farcaster.prom.$$;  then mv {{ NODE_EXPORTER_PATH }}/farcaster.prom.$$ {{ NODE_EXPORTER_PATH }}/farcaster.prom; else rm {{ NODE_EXPORTER_PATH }}/farcaster.prom.$$; fi 

